# deploy/bintrail-shim.service
[Unit]
Description=bintrail shim - BYOS time-travel SQL backend for ProxySQL
Documentation=https://github.com/dbtrail/bintrail/blob/main/docs/byos-time-travel-sql.md
After=network-online.target proxysql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/etc/bintrail
EnvironmentFile=/etc/bintrail/.bintrail.env
ExecStart=/usr/local/bin/bintrail shim --shim-config /etc/bintrail/shim.yaml
Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
